# Phidget/phidget.service
[Unit]
Description=Phidget Network Server
[Service]
Type=simple
ExecStart=/usr/bin/phidget22networkserver
[Install]
WantedBy=multi-user.target
